# scripts/launchd_deadman.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.hikari.deadman</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/__OPERATOR__/.local/bin/uv</string>
    <string>run</string>
    <string>python</string>
    <string>scripts/dead_man.py</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/__OPERATOR__/agents/hikari-agent</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HIKARI_DEADMAN_BOT_TOKEN</key><string>__DEADMAN_TOKEN__</string>
    <key>OWNER_TELEGRAM_ID</key><string>__OWNER__</string>
  </dict>
  <key>StartInterval</key><integer>300</integer>
  <key>StandardOutPath</key><string>/Users/__OPERATOR__/agents/hikari-agent/data/logs/deadman.out.log</string>
  <key>StandardErrorPath</key><string>/Users/__OPERATOR__/agents/hikari-agent/data/logs/deadman.err.log</string>
</dict>
</plist>
